# deploy/launchd/com.iterwheel.voyager.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.iterwheel.voyager.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>set -a &amp;&amp; source /Users/frank/.voyager/bridge.env &amp;&amp; set +a &amp;&amp; exec /Users/frank/.voyager/.venv/bin/vyg bridge serve --host 127.0.0.1 --port 8787</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/frank/Library/Logs/voyager/bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/frank/Library/Logs/voyager/bridge.err.log</string>

  <key>Umask</key>
  <integer>63</integer>
</dict>
</plist>
